# Ejemplos_practicos/Targets/Maquinas_quiosco/X11-apps/xeyes.service
[Unit]
Description=Añade cliente gráfico xeyes en segundo plano en /home/quiosco/.xinitrc
Documentation=https://github.com/adriisotuu/proyecto_maquina_quiosco/tree/master/Documentacion
After=xterm.service
Before=window_manager.service

[Service]
Type=oneshot
ExecStart=/bin/xeyes_xinitrc.sh

[Install]
WantedBy=quiosco_X11-apps.target
